# i3/99-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/vorakl/.Xauthority", RUN+="/home/vorakl/bin/set-monitors.sh"